# roles/rq/files/rq.target
[Unit]
Description=RQ Workers

[Install]
WantedBy=multi-user.target
